# t2/crontab
# m h  dom mon dow   command
* * * * * cd homework/sp/ICS/t2 ;./cr.sh
